# deploy/turing408.service
[Unit]
Description=Turing 408 FastAPI Service
After=network.target mysql.service

[Service]
WorkingDirectory=/opt/turing408/backend
ExecStart=/opt/turing408/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target

